# jarvis-voice.service
[Unit]
Description=JARVIS Voice Pipeline (STT + TTS + Orchestrator)
After=network.target jarvis-mkiii.service
Wants=jarvis-mkiii.service

[Service]
Type=simple
WorkingDirectory=/home/k/JARVIS-MKIII/backend
ExecStartPre=/bin/sleep 5
ExecStart=/home/k/JARVIS-MKIII/venv/bin/python3 voice/voice_orchestrator.py
Restart=on-failure
RestartSec=5
KillMode=control-group
Environment="PYTHONPATH=/home/k/JARVIS-MKIII/backend"

[Install]
WantedBy=default.target
